# systemd/oracle-postgres-backup.service
[Unit]
Description=Back up Oracle PostgreSQL
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/bots/oracle
ExecStart=/home/ubuntu/bots/oracle/scripts/backup-postgres
TimeoutStartSec=30min

