# systemd/oracle-brawlstars-claim.timer
[Unit]
Description=Claim Brawl Stars Store rewards once daily

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
RandomizedDelaySec=0
Persistent=true
AccuracySec=1s
Unit=oracle-brawlstars-claim.service

[Install]
WantedBy=timers.target
