# systemd/platform-brawlstars-claim.service
[Unit]
Description=Claim daily Brawl Stars Store rewards
After=docker.service network-online.target
Wants=network-online.target
StartLimitIntervalSec=3h
StartLimitBurst=3

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/srv/platform/infra
ExecStart=/usr/bin/flock -n /tmp/platform-brawlstars-claim.lock /srv/platform/infra/scripts/claim-brawlstars-reward
TimeoutStartSec=20min
Restart=on-failure
RestartSec=30min
RestartPreventExitStatus=2
