# systemd/sago-cloud-pr-media-prune.service
[Unit]
Description=Evict Sago Cloud PR media above its high watermark
RequiresMountsFor=/srv/pr-media

[Service]
Type=oneshot
Environment=GH_CONFIG_DIR=/home/ubuntu/.config/gh
ExecStart=/usr/local/libexec/sago-cloud-pr-media-prune
