# systemd/sago-cloud-health-watch.timer
[Unit]
Description=Check Sago Cloud container health every five minutes

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=30s
AccuracySec=30s
Persistent=true
Unit=sago-cloud-health-watch.service

[Install]
WantedBy=timers.target
