# systemd/sago-cloud-pr-media-prune.timer
[Unit]
Description=Check Sago Cloud PR media high watermark daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
