# systemd/sago-cloud-health-watch.service
[Unit]
Description=Restart unhealthy Sago Cloud containers with a cooldown
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/
Environment=OPERATIONS_ROOT=/srv/sago-cloud/operations
Environment=SAGO_CLOUD_ROOT=/srv/sago-cloud
ExecStart=/usr/local/libexec/sago-cloud-health-watch/current/scripts/restart-unhealthy-services
TimeoutStartSec=10min
