# systemd/sago-cloud-pr-media-verify.service
[Unit]
Description=Verify Sago Cloud PR media content hashes
RequiresMountsFor=/srv/pr-media

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/sago-cloud-pr-media-verify
